# deploy/dip-alert-deploy.service
[Unit]
Description=Dip Alert pull-based deploy (runs deploy.sh once)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dipalert
Group=dipalert
WorkingDirectory=/home/dipalert/dip-alert-app
# Deploy-only creds for the failure WhatsApp ping (DEPLOY_ALERT_PHONE/APIKEY).
# Deliberately NOT dip-alert.env: the app service must not carry these. Leading
# '-' = don't fail if the file is absent (alerts simply stay disabled). Re-read
# on every run, so editing the file needs no daemon-reload.
EnvironmentFile=-/etc/dip-alert/deploy-alert.env
ExecStart=/home/dipalert/dip-alert-app/deploy/deploy.sh
UMask=0077
# Leave enough time for a cold dependency install plus the regression gate.
TimeoutStartSec=30min
